# ops/dev-runtime/singrapw-dev-supabase.service
[Unit]
Description=SingraPW Dev Supabase Stack
After=network-online.target docker.service
Wants=network-online.target docker.service

[Service]
Type=oneshot
WorkingDirectory=/var/www/singrapw-dev
ExecStart=/usr/bin/npx supabase start
ExecStop=/usr/bin/npx supabase stop
RemainAfterExit=yes
TimeoutStartSec=600
TimeoutStopSec=120

[Install]
WantedBy=multi-user.target
